# scripts/tests.cron
*/5 * * * * /bin/bash /tests/scripts/start.sh
